# deploy/supervisor/oeparts-queue-worker.conf
; Supervisor config for the OeParts queue worker.
;
; Install: copy this file to /etc/supervisor/conf.d/, replace every
; /path/to/oeparts placeholder with the real deployment path, then run:
;   supervisorctl reread
;   supervisorctl update
;   supervisorctl start oeparts-queue-worker:*
;
; --queue=critical,default lists both queues in priority order — order
; confirmations, OTP, password resets, invoice generation, and Airwallex
; payment-webhook processing all run on 'critical'. Omitting it would
; silently process only 'default' and never touch 'critical'.
;
; stopwaitsecs must exceed the longest-running job's own $timeout
; (ProcessCsvImport.php declares $timeout=3600 for large catalog
; imports). Supervisor's own default (10s) would SIGKILL a legitimately
; still-running hour-long import the moment a deploy restarts the
; worker; 3700 gives it the full hour plus a buffer to finish cleanly.
;
; --max-time=3600 is the standard Laravel-recommended periodic self
; restart for long-running PHP worker processes (mitigates memory
; growth); autorestart=true respawns it immediately afterward.

[program:oeparts-queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/oeparts/artisan queue:work redis --queue=critical,default --sleep=3 --tries=3 --max-time=3600
directory=/path/to/oeparts
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
redirect_stderr=true
stdout_logfile=/path/to/oeparts/storage/logs/queue-worker.log
stopwaitsecs=3700
; user=www-data ; uncomment and set to your web server's user
